# deploy/systemd/agent-tool-firewall.service
[Unit]
Description=Agent Tool Firewall (default-deny policy gateway for LLM/agent tool calls)
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/agent-tool-firewall
Environment=BIND_ADDR=127.0.0.1:8475
Environment=POLICY_PATH=/etc/agent-tool-firewall/policy.yaml
Environment=AUDIT_LOG_PATH=/var/lib/agent-tool-firewall/audit.jsonl
Environment=SERVICE_TOKEN_PATH=/run/agent-tool-firewall/service-token

# Filesystem isolation
DynamicUser=yes
StateDirectory=agent-tool-firewall
LogsDirectory=agent-tool-firewall
ReadOnlyPaths=/etc/agent-tool-firewall
ReadWritePaths=/var/lib/agent-tool-firewall
PrivateNetwork=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectProc=invisible
ProcSubset=pid

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Privilege restriction
NoNewPrivileges=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictRealtime=yes

# Capability bounding
CapabilityBoundingSet=
AmbientCapabilities=

# Namespace restrictions
RestrictNamespaces=yes
RestrictAddressFamilies=AF_UNIX
PrivateUsers=yes

# Syscall filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @mount @clock @debug @swap @reboot @raw-io @module @cpu-emulation @obsolete
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM

# Resource limits
MemoryMax=256M
CPUQuota=25%
TasksMax=32
LimitNOFILE=512
LimitCORE=0

Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
